# rootfs/etc/services.d/swarm-client/run
#!/usr/bin/with-contenv bash

RESET='\E[0m'
YELLOW='\E[1;33m'

echo -e "${YELLOW}Starting swarm-client${RESET}"

s6-setuidgid jenkins /usr/bin/start-swarm-client
